# deploy/quant-monitor.service
[Unit]
Description=Quant Futures Testnet Live Monitor
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/quant-futures-bot
EnvironmentFile=/opt/quant-futures-bot/.env
ExecStart=/opt/quant-futures-bot/.venv/bin/python -m quant_futures_bot.live_monitor --poll-seconds 60
Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target

