# deploy/quant-altcoin-optimizer.service
[Unit]
Description=Quant Futures Altcoin Aggressive Rolling Backtest
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
WorkingDirectory=/opt/quant-futures-bot
EnvironmentFile=/opt/quant-futures-bot/.env
ExecStartPre=/opt/miniconda/envs/quant-bot/bin/python -m quant_futures_bot.historical_data --no-include-main --include-altcoin-top-volume --top 100 --start 2022-01-01 --end now --timeframes 2h --root quant_futures_bot/data/historical_ohlcv/binance_usdt_futures --sleep-seconds 0.02
ExecStart=/opt/miniconda/envs/quant-bot/bin/python -m quant_futures_bot.auto_altcoin_optimizer --run-once --top 100 --limit 30000 --timeframes 2h --show 30 --min-score 0 --max-leaders 30 --min-trades 8 --min-side-ratio 0 --fold-count 8 --min-profitable-fold-ratio 0.50 --max-hold-bars-2h 12 --extended-hold-bars-2h 6 --min-profit-to-extend 0.03 --trailing-after-max-hold-pct 0.03 --fee-rate 0.0002 --funding-cost-rate-per-8h 0.0001 --fetch-timeout-ms 15000 --fetch-retries 1 --strategy-workers 4 --data-root quant_futures_bot/data/historical_ohlcv/binance_usdt_futures --data-start 2022-01-01 --data-end now --output quant_futures_bot/data/altcoin_top100_2h_2022_backtest.csv --latest-json quant_futures_bot/data/altcoin_strategy_latest.json
TimeoutStartSec=12h
