# deploy/quant-altcoin-paper.timer
[Unit]
Description=Run Altcoin Paper Trading Status Every 15 Minutes

[Timer]
OnBootSec=5min
OnUnitActiveSec=15min
Unit=quant-altcoin-paper.service

[Install]
WantedBy=timers.target
